# webapp-base-iac/src/packer/webapp/conf/supervisor/webapp.conf
;exitcode 15 is SITGERM

[program:nginx]
command=nginx -g "daemon off;"
autorestart=unexpected
exitcodes=15

[program:php7-fpm]
command=php-fpm7\.2 --nodaemonize
autorestart=unexpected
exitcodes=15